# templates/systemd/kitsunebi@.service
[Unit]
Description=kitsunebi Minecraft instance %i
After=network-online.target
Wants=network-online.target

[Service]
User=minecraft
Group=minecraft
WorkingDirectory=/srv/kitsunebi/instances/%i/data
EnvironmentFile=/etc/kitsunebi/instances/%i.env
EnvironmentFile=-/etc/kitsunebi/secrets/%i.env
ExecStart=/usr/bin/java $JAVA_FLAGS -jar /srv/kitsunebi/instances/%i/data/server.jar nogui
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kitsunebi-%i

[Install]
WantedBy=kitsunebi.target
