# templates/systemd/kitsunebi.target
[Unit]
Description=kitsunebi Minecraft instance group

[Install]
WantedBy=multi-user.target
